# deploy/vigyan.service
[Unit]
Description=VigyanLLM Primer Design Backend
After=network.target postgresql.target
Wants=postgresql.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/vigyanpilot

# Load env vars from .env before the process starts
EnvironmentFile=/home/ubuntu/vigyanpilot/.env

# Activate venv and start gunicorn
ExecStart=/home/ubuntu/vigyanpilot/.venv/bin/gunicorn \
    -c deploy/gunicorn.conf.py \
    wsgi:app

# Restart on failure
Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
