# deploy/vigyan.service
[Unit]
Description=VigyanLLM Primer Design Backend
Documentation=https://vigyanllm.in/docs
After=network-online.target postgresql.service
Wants=postgresql.service
Requires=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/vigyanpilot

# ── Environment ────────────────────────────────────────────────────────────
Environment=PATH=/home/ubuntu/vigyanpilot/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=/home/ubuntu/vigyanpilot/.env

# ── Pre-start ──────────────────────────────────────────────────────────────
ExecStartPre=/home/ubuntu/vigyanpilot/.venv/bin/python \
    /home/ubuntu/vigyanpilot/deploy/migrations/migrate.py

# ── Application ────────────────────────────────────────────────────────────
ExecStart=/home/ubuntu/vigyanpilot/.venv/bin/gunicorn \
    -c /home/ubuntu/vigyanpilot/deploy/gunicorn.conf.py \
    wsgi:app

# ── Process Lifecycle ──────────────────────────────────────────────────────
Restart=on-failure
RestartSec=5
RestartSteps=3
RestartMaxDelaySec=30
TimeoutStartSec=60
TimeoutStopSec=30

# ── Resource Limits ────────────────────────────────────────────────────────
LimitNOFILE=65536
LimitNPROC=4096
LimitCORE=0

# ── Security Hardening ─────────────────────────────────────────────────────
# Moderate hardening — keeps the app functional while preventing common
# privilege escalation and information disclosure vectors.

# Prevent privilege escalation via setuid binaries
NoNewPrivileges=true

# Restrict filesystem access — writeable only where needed
ProtectSystem=strict
ReadWritePaths=/home/ubuntu/vigyanpilot /var/log/vigyan /tmp

# Isolate /tmp — noexec, nosuid, nodev
PrivateTmp=true

# Kernel protection (safe on all systems)
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectClock=true

# Isolate from hostname changes
ProtectHostname=true

# Remove inherited IPC resources on exit
RemoveIPC=true

# NOTE: ProtectHome=true is intentionally omitted — it blocks access to
# /home/ubuntu/vigyanpilot even with ReadWritePaths on some systemd versions.
# MemoryDenyWriteExecute is omitted — it breaks Python C extensions.
# SystemCallFilter/CapabilityBoundingSet are omitted — they interfere with
# gunicorn worker spawning and database driver socket connections.

[Install]
WantedBy=multi-user.target
